# goertzel/phreak.timer
[Timer]
OnBootSec=90

[Install]
WantedBy=timers.target
